# parser/crontab
# Run Scrapy every minute
* * * * * cd /app && python main.py
